# divemo.service
# session-divemo systemd unit (TEMPLATE — not installed by the repo).
#
# A system service (visible in `systemctl` to all admins) running as a
# dedicated, unprivileged `divemo` user. Suggested layout:
#
#   code:   /opt/session-divemo
#   config: /etc/divemo/config.yaml   (mode 640 root:divemo — holds the token)
#   state:  /var/lib/divemo/          (created automatically via StateDirectory)
#
# Install (see README for the full walkthrough):
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin divemo
#   sudo git clone <repo-url> /opt/session-divemo
#   sudo install -d -m 750 -o root -g divemo /etc/divemo
#   sudo install -m 640 -o root -g divemo config.sample.yaml /etc/divemo/config.yaml
#   sudoedit /etc/divemo/config.yaml   # fill in values; set state_db: /var/lib/divemo/divemo.db
#   sudo cp divemo.service /etc/systemd/system/divemo.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now divemo.service
#
# Logs:    journalctl -u divemo -f
# Status:  systemctl status divemo

[Unit]
Description=session-divemo — Discord version monitor for Session releases
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=divemo
Group=divemo
WorkingDirectory=/opt/session-divemo
ExecStart=/usr/bin/python3 /opt/session-divemo/bot.py -c /etc/divemo/config.yaml
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1 PYTHONDONTWRITEBYTECODE=1

# State lives in /var/lib/divemo (auto-created and chowned to the service user).
# Point state_db at it in the config: state_db: "/var/lib/divemo/divemo.db"
StateDirectory=divemo

# Hardening: the service only needs outbound network, read-only code/config,
# and its StateDirectory for writes.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
